# ops/systemd/chatgptrest-issue-graph-export.service
[Unit]
Description=ChatgptREST issue graph exporter
After=chatgptrest-api.service

[Service]
Type=oneshot
WorkingDirectory=__CHATGPTREST_ROOT__
EnvironmentFile=%h/.config/chatgptrest/chatgptrest.env
ExecStart=__CHATGPTREST_ROOT__/.venv/bin/python __CHATGPTREST_ROOT__/ops/export_issue_graph.py
