# ops/systemd/chatgptrest-issue-views-export.timer
[Unit]
Description=Periodic ChatgptREST issue view export

[Timer]
OnBootSec=5m
OnUnitActiveSec=15m
Unit=chatgptrest-issue-views-export.service

[Install]
WantedBy=timers.target
